# etc/supervisor/conf.d/vncd.conf
[program:websockify]
command=bash -c '/dockerstartup/d.sh'
autorestart=true
